# tools/oc-personal-runner/systemd/oc-personal.service
[Unit]
Description=oc-personal-runner — chat-completions front for Jarvis (Phase 8)
After=network-online.target
Wants=network-online.target
# Soft dependency on brain-sync — we want the vault freshly pulled before
# the agent answers, but it's fine if the timer fires after we start.
After=brain-sync.service

[Service]
Type=simple
# __RUN_USER__ and __PROJECT_ROOT__ are sentinel placeholders that deploy.sh
# substitutes at install time (see tools/oc-personal-runner/deploy.sh:cmd_install).
# This file is a template, not a directly-installable unit — the post-sed
# validator in deploy.sh refuses to install if any __FOO__ token survives.
User=__RUN_USER__
WorkingDirectory=__PROJECT_ROOT__/tools/oc-personal-runner
# Override the rest in a drop-in (/etc/systemd/system/oc-personal.service.d/override.conf)
# rather than editing this file. ANTHROPIC_API_KEY especially should not live
# in the unit text — use EnvironmentFile= pointing at a 0600-permissioned env file.
Environment=OC_LISTEN_HOST=0.0.0.0
Environment=OC_LISTEN_PORT=8080
Environment=OC_LMSTUDIO_URL=http://192.168.1.108:1234
Environment=OC_PERSONAL_MODEL=oc-personal
Environment=OC_ANTHROPIC_MODEL=claude-sonnet-4-6
Environment=OC_MAX_TURNS=4
Environment=OC_MAX_OUTPUT_TOKENS=200
Environment=BRAIN_VAULT_PATH=/srv/2ndbrain
# OC_MCP_SERVERS is the JSON catalog of stdio MCP children the runner
# spawns at startup. Each entry needs `command` and `args`; optional
# `env` overrides the spawn env, optional `stderr` redirects the child's
# stderr to that path (otherwise mcp.stdio_client silently discards it,
# which makes "Connection closed" errors essentially undebuggable).
#
# Quoted so systemd does not split the JSON on whitespace into multiple
# KEY=VAL pairs — same footgun as the old OC_BRAIN_MCP_ARGS line.
#
# Two-MCP setup: brain-mcp (wiki) + google-mcp (Gmail + Calendar). Drop
# the google entry and re-render with deploy.sh if you don't want it.
Environment="OC_MCP_SERVERS={\"brain\":{\"command\":\"__PROJECT_ROOT__/tools/brain-mcp/.venv/bin/python\",\"args\":[\"-m\",\"brain_mcp.server\"],\"env\":{\"BRAIN_VAULT_PATH\":\"/srv/2ndbrain\"},\"stderr\":\"/tmp/brain-mcp.err\"},\"google\":{\"command\":\"__PROJECT_ROOT__/tools/google-mcp/.venv/bin/python\",\"args\":[\"-m\",\"google_mcp.server\"],\"env\":{\"GOOGLE_CREDENTIALS_PATH\":\"__USER_HOME__/.config/oc-personal/google-credentials.json\",\"GOOGLE_TOKEN_PATH\":\"__USER_HOME__/.config/oc-personal/google-token.json\"},\"stderr\":\"/tmp/google-mcp.err\"}}"
EnvironmentFile=-/etc/oc-personal/secrets.env
ExecStart=__PROJECT_ROOT__/tools/oc-personal-runner/.venv/bin/python -m oc_personal.server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
